# share/system/systemd/NetRecorder-record.service
[Unit]
Description = NetRecorder record daemon
After = network-online.target

[Service]
Type = oneshot
ExecStart = /var/www/Net-Recorder/share/service/record.pl
#User = apache
#Group = apache

[Install]
WantedBy = multi-user.target
